# test.json
[
  {
    "id": "T1",
    "title": "Registry Run Key (Engineer WS)",
    "domain": "OT",
    "severity": "high",
    "technique": "T1547.001",
    "iec": ["SR 2.3", "SR 3.1"],
    "description": "Simulate persistence by logging a Run key creation on an engineering workstation. Uses Filebeat-tailed JSON and optional Winlogbeat/Sysmon.",
    "kql": "scenario:\"T1_registry_runkey\" OR (event.provider:\"Microsoft-Windows-Sysmon\" and winlog.event_data.TargetObject:*CurrentVersion\\\\Run*)",
    "sigma": "title: Registry Run Key (Sim)\\ndetection:\\n  selection:\\n    scenario: T1_registry_runkey\\n  condition: selection\\nlevel: high",
    "assets": ["eng-ws (Windows) with Filebeat", "optional: Sysmon + Winlogbeat"],
    "playbook": [
      "Append JSON to C:\\\\port254\\\\demo.log via Ansible win_shell.",
      "Filebeat forwards to Elasticsearch.",
      "Open Kibana Discover filtered by run_id."
    ],
    "notes": "Safe, reversible.",
    "video": "https://www.youtube.com/embed/UeLJwc8hwlM",
    "tags": ["persistence","windows","filebeat","safe"],
    "api_route": "/api/run/t1"
  },
  {
    "id": "T2",
    "title": "PowerShell IEX (Local)",
    "domain": "OT",
    "severity": "high",
    "technique": "T1059.001",
    "iec": ["SR 2.4","SR 3.1"],
    "description": "Detect encoded/command execution patterns typical of IEX payload delivery.",
    "kql": "process.image:*powershell.exe AND process.command_line:(*DownloadString* OR *IEX*)",
    "sigma": "title: Suspicious PowerShell IEX\\ndetection:\\n  selection:\\n    process.image: *powershell.exe\\n    process.command_line|contains: DownloadString\\n  condition: selection",
    "assets": ["eng-ws (Windows) with Filebeat/Winlogbeat"],
    "playbook": ["Append JSON to demo log with suspicious commandline."],
    "notes": "Simulated or safe benign script run.",
    "video": "",
    "tags": ["execution","windows","sysmon","filebeat"],
    "api_route": "/api/run/t2"
  },
  {
    "id": "T3",
    "title": "PLC Program Change (simulated upload)",
    "domain": "OT",
    "severity": "critical",
    "technique": "T0832",
    "iec": ["MR 2.2","SR 1.1"],
    "description": "Touch PLC project files and log uploads to simulate a program change.",
    "kql": "scenario:\"T3_plc_upload\" OR file.path:*demo_plc_project.proj",
    "sigma": "title: PLC Upload (Sim)\\ndetection:\\n  selection:\\n    scenario: T3_plc_upload\\n  condition: selection",
    "assets": ["plc1, plc-gateway, historian (Filebeat)"],
    "playbook": [
      "Copy demo_plc_project.proj to /opt/openplc/projects.",
      "Append gateway/historian event to /var/log/port254/demo.log."
    ],
    "notes": "Non-destructive, safe in lab.",
    "video": "",
    "tags": ["plc","upload","historian","filebeat"],
    "api_route": "/api/run/t3"
  },
  {
    "id": "T4",
    "title": "Modbus Write to Critical Register",
    "domain": "OT",
    "severity": "medium",
    "technique": "T0868",
    "iec": ["SR 3.2"],
    "description": "Detect Modbus/TCP writes to critical register addresses; pair network IDS log with historian note.",
    "kql": "scenario:\"T4_modbus_write\" OR (network.transport:tcp AND destination.port:502)",
    "sigma": "title: Modbus Write\\ndetection:\\n  selection:\\n    scenario: T4_modbus_write\\n  condition: selection",
    "assets": ["ids (Suricata EVE), historian (Filebeat)"],
    "playbook": ["Run pymodbus client from jump host to write a safe test value."],
    "notes": "Safe value; ensure isolated lab.",
    "video": "",
    "tags": ["network","modbus","ids","historian"],
    "api_route": "/api/run/t4"
  },
  {
    "id": "T5",
    "title": "Firmware File Replacement (server)",
    "domain": "OT",
    "severity": "high",
    "technique": "T0822",
    "iec": ["SR 1.1","SR 2.2"],
    "description": "Detect replacement of firmware binaries on the update server (checksum recorded).",
    "kql": "scenario:\"T5_firmware_replace\" OR event.action:\"firmware.replace\"",
    "sigma": "title: Firmware Replace\\ndetection:\\n  selection:\\n    scenario: T5_firmware_replace\\n  condition: selection",
    "assets": ["file-server (Filebeat)","plc1 (optional)"],
    "playbook": ["Copy demo_firmware_v2.bin and log sha256 to demo.log."],
    "notes": "Safe demo; no device flash.",
    "video": "",
    "tags": ["firmware","integrity","filebeat"],
    "api_route": "/api/run/t5"
  },
  {
    "id": "T6",
    "title": "NTP Time Tampering",
    "domain": "OT",
    "severity": "medium",
    "technique": "T0814",
    "iec": ["SR 3.3"],
    "description": "Simulate clock shift via log events on NTP/historian/eng-ws; detect time deltas.",
    "kql": "scenario:\"T6_ntp_shift\" OR event.action:\"ntp.adjust\"",
    "sigma": "title: NTP Adjust\\ndetection:\\n  selection:\\n    scenario: T6_ntp_shift\\n  condition: selection",
    "assets": ["ntp, historian, eng-ws (Filebeat)"],
    "playbook": ["Append ntp.adjust JSON to demo logs."],
    "notes": "Prefer simulated; real clock changes can disrupt services.",
    "video": "",
    "tags": ["time","ntp","filebeat"],
    "api_route": "/api/run/t6"
  },
  {
    "id": "F1",
    "title": "FinTech — High-Value Transaction Spike",
    "domain": "FinTech",
    "severity": "high",
    "technique": "T1071",
    "iec": [],
    "description": "Simulate a payment event > threshold to validate anomaly rules; useful for fraud runbooks.",
    "kql": "scenario:\"F1_high_value_txn\" OR (service:payments AND event.action:transaction.create AND amount >= 1000000)",
    "sigma": "title: High Value TXN\\ndetection:\\n  selection:\\n    scenario: F1_high_value_txn\\n  condition: selection",
    "assets": ["payments-app log (Filebeat)"],
    "playbook": ["Append JSON payment event with large amount and run_id."],
    "notes": "Synthetic, sanitized.",
    "video": "",
    "tags": ["fintech","fraud","filebeat"],
    "api_route": "/api/run/f1"
  },
  {
    "id": "F2",
    "title": "FinTech — Credential Stuffing Burst",
    "domain": "FinTech",
    "severity": "medium",
    "technique": "T1110",
    "iec": [],
    "description": "Burst of failed logins from single IP over short interval; validate rate-limit + detection.",
    "kql": "scenario:\"F2_credential_stuff\" OR (event.action:login.failed AND source.ip:* AND event.dataset:auth)",
    "sigma": "title: Credential Stuffing\\ndetection:\\n  selection:\\n    scenario: F2_credential_stuff\\n  condition: selection",
    "assets": ["auth service logs (Filebeat)"],
    "playbook": ["Append N failed login JSONs with same source.ip and run_id."],
    "notes": "No real users; sanitized.",
    "video": "",
    "tags": ["fintech","auth","bruteforce","filebeat"],
    "api_route": "/api/run/f2"
  }
]
